# systemd/port-peeker.service
# port-peeker - HTTP healthcheck agent
#
# Install (run as root):
#   install -m 755 bin/port-peeker-linux-arm64 /usr/local/bin/port-peeker
#   install -m 644 systemd/port-peeker.service /etc/systemd/system/port-peeker.service
#   systemctl daemon-reload
#   systemctl enable --now port-peeker
#   curl -s http://127.0.0.1:24365/healthz
#
# About process-name resolution (?process=NAME):
#   This sample runs as User=root so that /proc/<pid>/fd is readable for every
#   process and ?process=NAME works against any owner. If you do not need
#   ?process=NAME (port LISTEN check only), switch to the safer
#   DynamicUser=yes variant shown in the comments below.

[Unit]
Description=port-peeker - HTTP healthcheck agent
Documentation=https://github.com/kawaz/port-peeker
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/port-peeker --listen :24365
Restart=always
RestartSec=5

# Default: full features including ?process=NAME resolution
User=root

# --- Safer alternative: port LISTEN check only -------------------------------
# Remove the "User=root" line above and enable these two lines instead:
#User=
#DynamicUser=yes
# In that case ?process=NAME falls back to (none) for processes owned by
# other users.
# ----------------------------------------------------------------------------

# Sandbox hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
